# ramdisk/etc/inittab
::sysinit:/etc/init.d/rcS
tty0::askfirst:-/bin/sh
